# deploy/aura-proxy/systemd/aura-proxy-ingress.service
[Unit]
Description=AURA AIWire proxy ingress sidecar
Documentation=https://github.com/H-XX-D/AURA/blob/main/docs/aiwire_proxy.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Environment=AURA_PROXY_LISTEN_HOST=127.0.0.1
Environment=AURA_PROXY_LISTEN_PORT=9101
Environment=AURA_PROXY_EGRESS_HOST=192.168.0.10
Environment=AURA_PROXY_EGRESS_PORT=9102
Environment=AURA_PROXY_BACKEND=auto
Environment=AURA_PROXY_METRICS=/var/log/aura/aura-proxy-ingress.metrics.json
Environment=AURA_PROXY_REPLAY_LOG=/var/log/aura/aura-proxy-ingress.replay.jsonl
ExecStart=/usr/bin/env aura-proxy ingress --listen-host ${AURA_PROXY_LISTEN_HOST} --listen-port ${AURA_PROXY_LISTEN_PORT} --egress-host ${AURA_PROXY_EGRESS_HOST} --egress-port ${AURA_PROXY_EGRESS_PORT} --backend ${AURA_PROXY_BACKEND} --metrics-output ${AURA_PROXY_METRICS} --replay-log-output ${AURA_PROXY_REPLAY_LOG}
Restart=on-failure
RestartSec=2
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
